# roles/setup/templates/codecheck.conf.j2
description "codecheck test running server AKA whale"
# Environment variables used in github.com/givery-technology/whale-server
## Paper trail
env PAPERTRAIL_HOST={{ papertrail_host }}
## Docker
# env DOCKER_TIMEOUT=600 # Default is 10 minutes

start on startup
stop on runlevel [!2345]

script
  forever stopall
  swapoff -a
  [ -f /swapfile ] && rm /swapfile

  SYSTEM_MEM=1000
  DOCKER_MEM=256
  DOCKER_SWAP_MEM=1000
  export DOCKER_IMAGE_NAME=givery/codecheck
  export DOCKER_IMAGE_VERSION=latest
  export DOCKER_OPTIONS="--volumes-from='cache-master' --memory=${DOCKER_MEM}M --memory-swap=${DOCKER_SWAP_MEM}M --memory-swappiness='0'"
  export DOCKER_CONCURRENCY=$(echo "($(free -m | awk 'NR==2 {print $2}')-$SYSTEM_MEM)/$DOCKER_MEM" | bc)
  SWAP_SIZE=$(echo "$DOCKER_SWAP_MEM * $DOCKER_CONCURRENCY" | bc)
  fallocate -l ${SWAP_SIZE}M /swapfile
  chmod 600 /swapfile
  mkswap /swapfile
  swapon /swapfile
  df -h
  cat /proc/swaps
  cat /proc/meminfo
  forever start --uid "whale" -l /dev/null -m 1 --append {{ source_path }}/server.js
end script
